# agent/deploy/ipgate-agent.service
[Unit]
Description=ipgate agent — nftables 放行名单管理（default-drop）
Documentation=https://github.com/lion1991/ipgate
# 防火墙提前：在网卡配置前建立 default-drop，关闭暴露窗口（同 ufw/firewalld 做法）。
Wants=network-pre.target
Before=network-pre.target
After=local-fs.target

[Service]
Type=simple
ExecStart=/usr/local/bin/ipgate-agent --config /etc/ipgate/config.json run
Restart=on-failure
RestartSec=2
# 启动即应用 ruleset；失败（如 nft 报错）则退出非 0 → 重启重试。
TimeoutStartSec=30

# --- 加硬（root 下仍尽量收敛）---
NoNewPrivileges=yes
ProtectSystem=strict
# /etc 可写：开/关 SSH 密码登录要改 /etc/ssh/sshd_config；重置 root 密码 chpasswd 要写
# /etc/shadow 并在 /etc 下建锁/临时/备份文件。没有它，ProtectSystem=strict 会让这些写操作
# 在服务的 mount namespace 里报「只读文件系统」（手动 shell 改没事，服务改不了）。
ReadWritePaths=/var/lib/ipgate /etc
ProtectHome=yes
PrivateTmp=yes
ProtectControlGroups=yes
ProtectKernelLogs=yes
# nft 需要 netlink，default-drop 需操作内核防火墙，故地址族必须放行 AF_NETLINK。
RestrictAddressFamilies=AF_UNIX AF_INET AF_INET6 AF_NETLINK

[Install]
WantedBy=multi-user.target
